# gunicorn.conf
[program:gunicorn]
command=/envs/lake-gallery/bin/gunicorn lakegallery.wsgi -b 0.0.0.0:8000
directory=/lakegallery